# live_build_config/includes.chroot_after_packages/etc/systemd/system/logrotate.timer
[Unit]
Description=Timer: Logs rotation

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=logrotate.service

[Install]
WantedBy=timers.target
