# live_build_config/includes.chroot_after_packages/usr/lib/systemd/system/cephos-host-ready.service
[Unit]
Description=Play startup melody when system is up
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/startup-melody.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
