# live_build_config/includes.chroot_after_packages/usr/lib/systemd/system/cephos-mds-manager.service
[Unit]
Description=Automatic Ceph MDS scaling agent
After=network-online.target ceph-mon.target
Wants=ceph-mon.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/cephos-mds-manager.sh
User=root
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
